# config/systemd/openclaw-gateway.service.d/env.conf
# Drop-in for openclaw-gateway.service — keeps secrets OUT of the unit file.
# Install to: ~/.config/systemd/user/openclaw-gateway.service.d/env.conf
#
# The referenced env file holds provider API keys, one per line:
#   OPENROUTER_API_KEY=sk-or-v1-...
# Create it with mode 600 and OUTSIDE this repo:
#   install -m 600 /dev/null ~/.openclaw/secrets/openrouter.env
#   printf 'OPENROUTER_API_KEY=%s\n' "$KEY" > ~/.openclaw/secrets/openrouter.env
#
# openclaw.json reads OPENROUTER_API_KEY from the process env (see
# models.providers.openrouter.apiKey -> { source: env, id: OPENROUTER_API_KEY }).

[Service]
EnvironmentFile=/home/USER/.openclaw/secrets/openrouter.env
